# supervisor.conf
[program:education_system]
user=education_system
stdout_logfile=/var/log/suervisor_stdout.log
redirect_stderr=true
command=/home/app/gunicorn_start.sh
